# supervisord.conf
[supervisord]
nodaemon=true

[program:front]
command=node /home/front/.next/standalone/server.js
autorestart=true
directory=/home/front
